# systemd/forex-ml-dashboard.service
[Unit]
Description=Forex ML Trading Web Dashboard
After=network.target forex-ml-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/forex-ml-trading
Environment="PATH=/home/forex-ml-trading/venv_pro/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/forex-ml-trading/venv_pro/bin/python /home/forex-ml-trading/web_dashboard/app.py
Restart=always
RestartSec=30
StandardOutput=append:/home/forex-ml-trading/logs/dashboard.log
StandardError=append:/home/forex-ml-trading/logs/dashboard_error.log

# حماية من إعادة التشغيل المفرطة
StartLimitInterval=600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target